# system_files/etc/profile.d/glorpfetch.sh
alias neofetch='/usr/bin/glorpfetch'
alias neowofetch='/usr/bin/glorpfetch'
alias fastfetch='/usr/bin/glorpfetch'
